# rpi-spotifycontrol.service.template
[Unit]
Description=RPi Spotify Control
After=network-online.target

[Service]
ExecStart=/usr/bin/python3 -u main.py
WorkingDirectory=$HOME/rpi-spotifycontrol/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=$USER

[Install]
WantedBy=network-online.target